# multimedia/autostart/detect-dvd.desktop
[Desktop Entry]
Type=Application
Name=Karoshi video DVD detection
Comment=Detect video DVDs and generate icon
Exec=karoshi-detect-dvd
